# s6-overlay/s6-rc.d/intervalhandler/run
#!/command/execlineb -P
/opt/hide.me/intervalhandler.py
